# tests/fixtures/systemd/benign.service
[Unit]
Description=Cleanup task

[Service]
ExecStart=/usr/bin/true
