# tests/fixtures/systemd/suspicious.service
[Unit]
Description=Telemetry updater

[Service]
ExecStart=/tmp/.cache/persist.sh
Environment=LD_PRELOAD=/tmp/.evil/libhook.so
Restart=always
